# core/chef/cookbooks/provisioner/templates/default/sws-init.erb
#!/bin/bash
# chkconfig:         2345 90 10
# description:       Service runner for sws
### BEGIN INIT INFO
# Provides:          sws
# Required-Start:    $syslog $network
# Required-Stop:     $syslog $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start the Simple Web Server
# Description:       Manages the Simple Web Server

start() {
    status && exit 0
    /usr/local/bin/sws --listen=:<%= @port %> --site=<%= @docroot %> </dev/null 2>&1 1>/var/log/sws.log &
    disown %1
}

stop() {
    if status; then
        pkill /usr/local/bin/sws
    else
        echo goiardi is not running
        exit 0
    fi
}

restart() {
    stop; start
}

status() {
    pgrep /usr/local/bin/sws
}

case "$1" in
    start|stop|status|restart) "$1";;
    *)
        echo "Usage: $0 {start|stop|status|restart}"
        exit 1;;
esac
