# core/chef/cookbooks/provisioner/templates/default/tftp.service.erb
[Unit]
Description=Rebar TFTP Server

[Service]
ExecStart=/usr/sbin/in.tftpd -s <%= @tftproot %>
StandardInput=socket