# chef/cookbooks/rebar/templates/default/dns-mgmt-init.erb
#!/bin/bash
# chkconfig:         2345 90 10
# description:       Service runner for rebar-dns-mgmt
### BEGIN INIT INFO
# Provides:          rebar-dns-mgmt
# Required-Start:    $syslog $network
# Required-Stop:     $syslog $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start the DNS Mgmt server
# Description:       Manages the DNS Mgmt server

start() {
    status && exit 0
    /usr/local/bin/rebar-dns-mgmt </dev/null 1>/var/log/rebar-dns-mgmt.log 2>&1 &
    disown %1
}

stop() {
    if status; then
        pkill rebar-dns-mgmt
    else
        echo rebar-dns-mgmt is not running
        exit 0
    fi
}

restart() {
    stop; start
}

status() {
    pgrep rebar-dns-mgmt
}

case "$1" in
    start|stop|status|restart) "$1";;
    *)
        echo "Usage: $0 {start|stop|status|restart}"
        exit 1;;
esac
